# daemon/gdocs.mount.in
[Mount]
Type=gdocs
Exec=@libexecdir@/gvfsd-gdocs
AutoMount=true
DBusName=org.gtk.vfs.mountpoint.gdocs
